# private_Library/private_LaunchAgents/io.github.reidond.bwssh-agent.plist.tmpl
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>io.github.reidond.bwssh-agent</string>

    <!-- launchd does NOT inherit the shell PATH, so we wrap in a login shell
         that activates mise — this resolves bwssh-agentd (and the bw / bitwarden
         CLI it shells out to) exactly like an interactive terminal, and survives
         mise version upgrades without a brittle hardcoded install path. `exec`
         replaces the shell so launchd supervises bwssh-agentd directly and sees
         its real exit status.

         KeepAlive uses SuccessfulExit=false (not plain true): restart only when
         the daemon crashes with a non-zero exit — e.g. the vaultwarden host
         unreachable at boot (EHOSTUNREACH in the logs) — and leave it alone on a
         clean exit. ThrottleInterval=20 backs off restarts during that boot-time
         LAN/network window so launchd doesn't hammer the unreachable host.

         Config (vaultwarden URL etc.) lives in ~/.config/bwssh and the bw CLI's
         own state, so no BW_* env vars are needed — same as the previous plist,
         which only injected PATH and ran fine. -->
    <key>ProgramArguments</key>
    <array>
        <string>/bin/zsh</string>
        <string>-lc</string>
        <string>eval "$(/opt/homebrew/bin/mise activate zsh 2&gt;/dev/null)"; exec bwssh-agentd --foreground</string>
    </array>

    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key>
        <string>{{ .chezmoi.homeDir }}</string>
    </dict>

    <key>WorkingDirectory</key>
    <string>{{ .chezmoi.homeDir }}</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>

    <!-- Back off between restart attempts (e.g. during a vaultwarden/LAN outage). -->
    <key>ThrottleInterval</key>
    <integer>20</integer>

    <key>ProcessType</key>
    <string>Background</string>

    <key>StandardOutPath</key>
    <string>/tmp/bwssh-agent.stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/bwssh-agent.stderr.log</string>
</dict>
</plist>
